# packages/fujin-cli/fujin_cli-0.16.0.tar.gz/fujin_cli-0.16.0/src/fujin/templates/default.service.j2
# All options are documented here https://www.freedesktop.org/software/systemd/man/latest/systemd.exec.html
[Unit]
Description={{ process_name }}

[Service]
User={{ user }}
Group={{ user }}
WorkingDirectory={{ app_dir }}
ExecStart={{ app_dir }}/{{ command }}
EnvironmentFile={{ app_dir }}/.env
Restart={% if process.timer %}no{% else %}always{% endif %}

[Install]
WantedBy=multi-user.target
