# packages/fujin-cli/fujin_cli-0.16.0.tar.gz/fujin_cli-0.16.0/src/fujin/templates/default.timer.j2
[Unit]
Description=Timer for {{ process_name }}

[Timer]
{% if process.timer.on_calendar -%}
OnCalendar={{ process.timer.on_calendar }}
{% endif -%}
{% if process.timer.on_boot_sec -%}
OnBootSec={{ process.timer.on_boot_sec }}
{% endif -%}
{% if process.timer.on_unit_active_sec -%}
OnUnitActiveSec={{ process.timer.on_unit_active_sec }}
{% endif -%}
{% if process.timer.on_active_sec -%}
OnActiveSec={{ process.timer.on_active_sec }}
{% endif -%}
{% if process.timer.randomized_delay_sec -%}
RandomizedDelaySec={{ process.timer.randomized_delay_sec }}
{% endif -%}
{% if process.timer.accuracy_sec -%}
AccuracySec={{ process.timer.accuracy_sec }}
{% endif -%}
Unit={{ process_name }}.service
Persistent={{ process.timer.persistent | lower }}

[Install]
WantedBy=timers.target
