# packages/keybird/keybird-1.0.12.tar.gz/keybird-1.0.12/keybird/systemd/pi-hid-bridge.service
[Unit]
Description=Pi HID Bridge Web App
After=network-online.target hid-gadget.service
Wants=network-online.target
Requires=hid-gadget.service

[Service]
Type=simple
User=root
WorkingDirectory=/home/pi/pi-hid-bridge
ExecStart=/usr/bin/python3 /home/pi/pi-hid-bridge/pi_kb.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target

