# packages/fujin-cli/fujin_cli-0.16.0-py3-none-any.whl/fujin/templates/web.service.j2
# All options are documented here https://www.freedesktop.org/software/systemd/man/latest/systemd.exec.html
# Inspiration was taken from here https://docs.gunicorn.org/en/stable/deploy.html#systemd
[Unit]
Description={{ process_name }}
{%- if process.socket %}
Requires={{ app_name }}.socket
{%- endif %}
After=network.target

[Service]
User={{ user }}
Group={{ user }}
RuntimeDirectory={{ app_name }}
WorkingDirectory={{ app_dir }}
ExecStart={{ app_dir }}/{{ command }}
EnvironmentFile={{ app_dir }}/.env
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
# if your app does not need administrative capabilities, let systemd know
ProtectSystem=strict

[Install]
WantedBy=multi-user.target
